# deployment/backup.service
[Unit]
Description=Backup Utah Campaign Finance Database
After=network.target

[Service]
Type=oneshot
User=root
Group=root

# Run the backup script
ExecStart=/var/www/polstats/deployment/backup.sh

# Logging
StandardOutput=append:/var/log/polstats/backup.log
StandardError=append:/var/log/polstats/backup-error.log

[Install]
WantedBy=multi-user.target
